# arch/save/archI/wol.service
[Unit]
Description=wol to trends

[Service]
Type=oneshot
ExecStart=wol a4:bb:6d:06:59:17
KillMode=process

[Install]
WantedBy=default.target
